# supervisord.conf
[supervisord]
nodaemon=true

[program:mumble]
command=sudo -u mumble /usr/sbin/murmurd -fg -ini /mumble-data/mumble-server.ini
autorestart=false
